# rc.local
#!/bin/bash
consul agent -ui -dev -client="0.0.0.0" > /dev/null 2> /dev/null &
exit 0
